# supervisor/xxai-iaa.ini
[program:xxai-iaa]
command=$EXE
loglevel=warn
autostart=true
autorestart=true
startretries=10
stopwaitsecs=1
startsecs=3
stdout_logfile=/var/log/supervisor/xxai-iaa.log
stderr_logfile=/var/log/supervisor/xxai-iaa.err.log
